# local_deploy/com.tychos.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tychos.worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>__TYCHOS_DIR__/tychos_skyfield/.venv/bin/python</string>
        <string>-m</string>
        <string>server.worker</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__TYCHOS_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>tychos_skyfield:tests:.</string>
        <key>PATH</key>
        <string>__TYCHOS_DIR__/tychos_skyfield/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__TYCHOS_DIR__/logs/worker.log</string>

    <key>StandardErrorPath</key>
    <string>__TYCHOS_DIR__/logs/worker.err.log</string>
</dict>
</plist>
